# server/reboot-notify/reboot-notify.service
[Unit]
Description=Send Slack notification on reboot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/reboot-notify.sh

[Install]
WantedBy=multi-user.target