# server/gpu-monitoring/idle-check/idle_gpu_killer.timer
[Unit]
Description=Run GPU Idle Killer at 22:00 KST daily

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
