# server/gpu-monitoring/idle-check/idle_gpu_killer.service
[Unit]
Description=GPU Idle Killer Script

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/idle_gpu_killer.sh
User=root
